# sisters-heartbeat.service
[Unit]
Description=Rank-42 Manifold Heartbeat API
After=network.target

[Service]
User=nous
WorkingDirectory=/home/nous/Aether
ExecStart=/home/nous/aether_env/bin/python3 /home/nous/Aether/heartbeat_api.py
Restart=always

[Install]
WantedBy=multi-user.target
